# jenkins_buildnode/rc.local
#!/bin/sh

GITEA_URL="https://gitea.tooling.${DOMAIN_NAME_SL}.${DOMAIN_NAME_TL}:3000/${ORG_NAME}/"
GITHUB_URL="https://github.com/Devoteam/"

if wget --spider "$GITEA_URL" >/dev/null 2>&1; then
    echo "Setting git to Gitea." >> /home/jenkins/install-log.txt
    export GIT_URL="$GITEA_URL"
else
    echo "Setting git to Github." >> /home/jenkins/install-log.txt
    export GIT_URL="$GITHUB_URL"
fi

echo "export GIT_URL=$GIT_URL" >> /etc/profile

/usr/bin/dockerd &
